# 02.busybox/rootfs/etc/inittab
# /etc/inittab: init(8) configuration.

# The default runlevel.
id:5:initdefault:

# Boot-time system configuration/initialization script.
# This is run first except when booting in emergency (-b) mode.
si::sysinit:/etc/init.d/rcS

l5:5:wait:/etc/rc.d/rc 5

#ttyO0:5:respawn:/bin/bash

z6:6:respawn:/sbin/sulogin
S:2345:respawn:/sbin/getty 115200 ttyO0
# /sbin/getty invocations for the runlevels.
#
# The "id" field MUST be the same as the last
# characters of the device (after "tty").
#
# Format:
#  <id>:<runlevels>:<action>:<process>
#

1:2345:respawn:/sbin/getty 38400 tty1

